# supervisord.conf
[program:telegram_bot]
command=python bot.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/bot.log
stderr_logfile=/var/log/bot_error.log
environment=PYTHONPATH='/app'